# files/profile.d/kubernetes.sh
#!/bin/sh
#
# Inclusion of Kubernetes Binaries into Path
export PATH=$PATH:/opt/kubernetes/server/bin