# templates/default/etc/systemd/system/kube-proxy.service.erb
[Unit]
Description=Kubernetes Proxy
Documentation=https://github.com/kubernetes/kubernetes/
After=network.target
Requires=<% @etcd_name %>.service
After=<% @etcd_name %>.service

[Service]
EnvironmentFile=/etc/default/kube-proxy
ExecStart=/usr/local/bin/hyperkube proxy $KUBE_PROXY_OPTS
Restart=on-failure
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target