# templates/default/etc/systemd/system/kubelet.service.erb
[Unit]
Description=Kubernetes Kubelet Service
Documentation=https://github.com/kubernetes/kubernetes/
After=network.target
Requires=<% @etcd_name %>.service
After=<% @etcd_name %>.service

[Service]
EnvironmentFile=/etc/default/kubelet
ExecStart=/usr/local/bin/hyperkube kubelet $KUBELET_OPTS
Restart=on-failure
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target