# templates/default/etc/systemd/system/kube-controller-manager.service.erb
[Unit]
Description=Kubernetes Controller Manager
Documentation=https://github.com/kubernetes/kubernetes/
After=network.target
Requires=<% @etcd_name %>.service
After=weave.service

[Service]
EnvironmentFile=/etc/default/kube-controller-manager
ExecStart=/opt/kubernetes/server/bin/kube-controller-manager $KUBE_CONTROLLER_MANAGER_OPTS
Restart=on-failure
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target