# templates/default/etc/systemd/system/kube-controller-manager.service.erb
[Unit]
Description=Kubernetes Controller Manager
Documentation=https://github.com/kubernetes/kubernetes/
After=network.target
Requires=<% @etcd_name %>.service
After=<% @etcd_name %>.service

[Service]
EnvironmentFile=/etc/default/kube-controller-manager
ExecStart=/usr/local/bin/hyperkube controller-manager $KUBE_CONTROLLER_MANAGER_OPTS
Restart=on-failure
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target