# templates/default/etc/systemd/system/kube-scheduler.service.erb
[Unit]
Description=Kubernetes Scheduler
Documentation=https://github.com/kubernetes/kubernetes/
After=network.target
Requires=<% @etcd_name %>.service
After=<% @etcd_name %>.service

[Service]
EnvironmentFile=/etc/default/kube-scheduler
ExecStart=/usr/local/bin/hyperkube scheduler $KUBE_SCHEDULER_OPTS
Restart=on-failure
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target